# sys-libs/conjuration/files/9999-initscript
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/dev-vcs/git/files/git-daemon.initd,v 1.1 2010/03/17 15:13:27 sping Exp $

##FIXME: Improve the socket access control mechanism.

pidfile="/run/conjuration-daemon.pid"
socket="${socket:-/run/conjuration.sock)"
socketowner="${socketowner:-root:root}"
socketperms="${socketperms:-ug=wr,o=}"

depend() {
	:
#        need net
#        use logger
}

start() {
	tempdir="$(mktemp -d)"
	chown "${socketowner}" "${tempdir}"
	chmod "${socketperms}" "${tempdir}"
	ebegin "Starting conjuration daemon"
	start-stop-daemon --start --quiet --background \
		--pidfile ${pidfile} --make-pidfile \
		--exec /usr/sbin/conjuration_daemon -- \
			"${socket}" \
			"${tempdir}"
	eend $?
}

stop() {
	ebegin "Stopping conjuration daemon"
	start-stop-daemon --stop --quiet \
		--pidfile ${pidfile}
	RET=$?
	##FIXME: Make the daemon actually respond to signals and shutdown properly, and remove this.
	rm "${socket}"
	eend ${RET}
}
